# systemd/.config/systemd/user/pbs-desktop-backup.timer
[Unit]
Description=Daily desktop backup to Proxmox Backup Server (user)

[Timer]
OnCalendar=*-*-* 08:00
Persistent=true

[Install]
WantedBy=timers.target
