# systemd/.config/systemd/user/pbs-desktop-backup.service
[Unit]
Description=Proxmox Backup - Desktop backup (user)

[Service]
Type=oneshot
Environment=PBS_REPOSITORY=austin-desktop@pbs!austin-desktop@pbs.theblacklodge.org:workstations
Environment=PBS_PASSWORD_FILE=/etc/proxmox-backup/austin-desktop.token
ExecStart=%h/.local/bin/pbs-desktop-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
