# systemd/.config/systemd/user/protonmail-bridge.service

[Unit]
Description=Proton Mail Bridge (headless CLI)

[Service]
Type=simple
Environment=QT_QPA_PLATFORM=offscreen
ExecStart=/usr/bin/protonmail-bridge --cli --noninteractive
Restart=on-failure

[Install]
WantedBy=default.target
